# english-corner-ai.service
[Unit]
Description=English Corner AI Backend Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ef
Group=ef
WorkingDirectory=/home/ef/english-corner-ai-backend
Environment=PATH=/home/ef/english-corner-ai-backend/venv/bin
ExecStart=/home/ef/english-corner-ai-backend/venv/bin/uvicorn rag_backend:app --host 0.0.0.0 --port 8443 --workers 1
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=english-corner-ai

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/home/ef/english-corner-ai-backend
ProtectHome=yes

[Install]
WantedBy=multi-user.target
