# english-corner-ai.service
[Unit]
Description=English Corner AI Backend Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ef
Group=ef
WorkingDirectory=/home/ef/english-corner-ai-backend
Environment=PATH=/home/ef/english-corner-ai-backend/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/ef/english-corner-ai-backend/venv/bin/uvicorn rag_backend:app --host 127.0.0.1 --port 8000 --workers 1
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=english-corner-ai

# Simplified security settings for compatibility
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
